# monitor/supervisor.conf

[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/crond -n
autostart=true
autorestart=true

